# src/modules/100-portron/etc/profile.d/globals
#!/bin/sh

# Export variables:
export \
scripts=/opt/scripts \
DISPLAY=:0
