# src/modules/100-portron/etc/profile.d/globals.sh
#!/bin/sh

# Export variables:
export \
DISPLAY=:0
