# src/modules/99-xterm/etc/profile.d/terminal.sh
#!/bin/sh

# Bring the cursor back if disabled by bootsplash:
echo -e "\e[?25h"
